# profile.d/phpstorm.sh
alias phpstorm='open -na "PhpStorm.app"'
