# profile.d/phpstorm.sh
phpstorm() {
    open -na "PhpStorm.app" --args "$@"
}
